# aggregator/crontab
# Scrapers
# Run daily at 3:00 AM
# 0 3 * * * /usr/local/bin/scrapers >> /var/log/scrapers.log 2>&1

# Run hourly
0 * * * * /usr/local/bin/scrapers >> /var/log/scrapers.log 2>&1

# Add more jobs below (email, matching, etc.)
